# debian/newscache.cron.daily
#!/bin/sh -e

# If it's configured
if [ -f /etc/newscache.conf ]; then
    /usr/sbin/newscacheclean > /dev/null
fi
